# deploy/templates/supervisor-hexugo.conf
; /etc/supervisor/conf.d/hexugo.conf
; Laravel queue workers + scheduler

[program:hexugo-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hexugo/current/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --backoff=5 --queue=high,default,low
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/hexugo/shared/storage/logs/queue-worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stopwaitsecs=3600

[program:hexugo-schedule]
command=/bin/bash -c "while true; do php /var/www/hexugo/current/artisan schedule:run --no-interaction >> /var/www/hexugo/shared/storage/logs/scheduler.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=deploy
redirect_stderr=true
stdout_logfile=/var/www/hexugo/shared/storage/logs/scheduler-supervisor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
